# lib/diffhound-webhook.service
[Unit]
Description=Diffhound PR Review Webhook Server
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/diffhound/lib
Environment="WEBHOOK_SECRET=CHANGEME"
Environment="DIFFHOUND_BIN=/home/ubuntu/diffhound/bin/diffhound"
Environment="REVIEW_REPO_PATH=/home/ubuntu/monorepo"
Environment="REVIEW_LOGIN=diffhound-bot"
Environment="DIFFHOUND_LOG_DIR=/home/ubuntu/logs/pr-reviews"
Environment="DIFFHOUND_MAX_CONCURRENT=2"
Environment="DIFFHOUND_DEBOUNCE_SECONDS=60"
Environment="HOME=/home/ubuntu"
Environment="PATH=/home/ubuntu/.local/bin:/home/ubuntu/.claude/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/gunicorn -w 1 -b 0.0.0.0:8090 --timeout 600 webhook:app
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=diffhound-webhook

[Install]
WantedBy=multi-user.target
